# systemd/palantir-web.service
[Unit]
Description=Palantir Web Service
After=redis.service
Requires=redis.service

[Service]
Type=simple
User=palantir
Group=palantir
WorkingDirectory=/opt/palantir
Environment=PALANTIR_ENV=production
EnvironmentFile=/opt/palantir/.env
ExecStart=/opt/palantir/.venv/bin/palantir-web
Restart=always
RestartSec=5
Nice=10
StandardOutput=journal
StandardError=journal

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/var/lib/palantir

[Install]
WantedBy=multi-user.target
