# systemd/palantir-vision.service
[Unit]
Description=Palantir Vision Service
After=redis.service
Requires=redis.service

[Service]
Type=simple
User=palantir
Group=palantir
WorkingDirectory=/opt/palantir
Environment=PALANTIR_ENV=production
EnvironmentFile=/opt/palantir/.env
ExecStart=/opt/palantir/.venv/bin/palantir-vision
Restart=always
RestartSec=5
Nice=-5
StandardOutput=journal
StandardError=journal

NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/palantir

# Camera device access
SupplementaryGroups=video

[Install]
WantedBy=multi-user.target
