# systemd/palintir-backup.timer
[Unit]
Description=Nightly Palintir database backup
Requires=palintir-backup.service

[Timer]
# Runs daily at 03:15 local time. Persistent=true catches up missed runs
# after reboots or planned downtime.
OnCalendar=*-*-* 03:15:00
Persistent=true
RandomizedDelaySec=15min
Unit=palintir-backup.service

[Install]
WantedBy=timers.target
